# infra/systemd/alliances-queue.service
# Install as a per-user unit (no sudo needed at deploy time):
#   ssh deploy@host
#   mkdir -p ~/.config/systemd/user
#   cp /var/www/alliances/backend/current/infra/systemd/alliances-queue.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alliances-queue
# Root must run `loginctl enable-linger deploy` once so the unit survives logout.
[Unit]
Description=Alliances Laravel queue worker
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/alliances/backend/current
ExecStart=/usr/bin/php artisan queue:work --queue=default --sleep=3 --tries=3 --backoff=30 --max-time=3600
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
