# infra/systemd/alliances-frontend.service
# Install as a per-user unit (no sudo needed at deploy time):
#   ssh deploy@host
#   mkdir -p ~/.config/systemd/user
#   cp /var/www/alliances/frontend/current/infra/systemd/alliances-frontend.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now alliances-frontend
# Root must run `loginctl enable-linger deploy` once so the unit survives logout.
[Unit]
Description=Alliances marketing site (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=/var/www/alliances/frontend/current
EnvironmentFile=/var/www/alliances/frontend/current/.env.production
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
