# infra/systemd/alliances-frontend.service
# Install: sudo cp infra/systemd/alliances-frontend.service /etc/systemd/system/
#          sudo systemctl daemon-reload
#          sudo systemctl enable --now alliances-frontend
[Unit]
Description=Alliances marketing site (Next.js)
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/var/www/alliances/frontend/current
EnvironmentFile=/var/www/alliances/frontend/current/.env.production
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=127.0.0.1
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
